# .config/systemd/user/default.target.wants/kmonad-lap.service
[Unit]
Description=KMonad configuration for Laptop Keyboard
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kmonad /home/alunity/scripts/lap-keyboard.kbd
Restart=on-failure

[Install]
WantedBy=default.target
